# timer-task.timer
[Unit]
Description=Timer

[Timer]
OnCalendar=:*0/15
Persistent=true

[Install]
MantedBy=timers.target